# root/etc/services.d/avahi/run
#!/usr/bin/with-contenv bash

exec avahi-daemon --no-drop-root
