# root/etc/services.d/mdns/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /usr/local/bin/docker_mdns "$NETWORK_INTERFACE"
